# crontab
*/5 * * * * php /home/milan/www/net/check.php
